# provision/ansible/roles/docker/templates/docker.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H unix:///var/run/docker.sock -H tcp://{{ VM.network.private.address }}:2375 {{ DOCKER.opts }}
