# provision/ansible/roles/clitools/files/crontab
01 * * * *  root [ -x /usr/local/bin/ct ] && /usr/local/bin/ct system:crontask &> /dev/null
#eof

